# deploy/modeld.service
[Unit]
Description=modeld-go (model manager)
After=network-online.target
Wants=network-online.target

[Service]
ExecStart=%h/Projects/modeld-go-1/bin/modeld --config %h/Projects/modeld-go-1/configs/models.yaml --addr :8080
WorkingDirectory=%h/Projects/modeld-go-1
Restart=on-failure
Environment=GOMAXPROCS=0
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
